# cron/crontab
# Chargement des variables d'environnement
# Exécution de la commande principale toutes les minutes
* * * * * /etc/launch_crontab.sh
